# deploy/gas-city-dashboard.service
[Unit]
Description=gas-city-dashboard — editorial ambient dashboard for a single gc operator (localhost-only, port 8082)
Documentation=file://%h/gas-city-dashboard/specs/architecture/overview.md
After=network-online.target
Wants=network-online.target

# Template via systemd's %h (user home) so the unit works on any operator's host
# when installed under `systemctl --user`. Override the env block below if your
# repo lives somewhere other than ~/gas-city-dashboard.

[Service]
Type=simple
WorkingDirectory=%h/gas-city-dashboard
Environment=NODE_ENV=production
Environment=PORT=8082
Environment=GC_SUPERVISOR_URL=http://127.0.0.1:8372
Environment=GC_CITY_NAME=gas-city
# GC_CITY_PATH pins the city root so the `gc bd` CLI write path
# (claim/close/nudge) doesn't depend on WorkingDirectory, which is the
# dashboard repo, not a city directory. Without it, writes fail with
# "not in a city directory". Override to match your actual city location.
Environment=GC_CITY_PATH=%h/gas-city
Environment=ADMIN_AUDIT_LOG_PATH=%h/.gc/events.jsonl
Environment=ADMIN_FRONTEND_DIST=%h/gas-city-dashboard/frontend/dist
# Refuse to start if a stale process holds the port. Bound to 127.0.0.1 only,
# so the ss check is sufficient — no public-interface race.
ExecStartPre=/bin/sh -c '! ss -tln "sport = :8082" | grep -q LISTEN || (echo "port 8082 already bound" >&2; exit 1)'
ExecStart=/usr/bin/node %h/gas-city-dashboard/backend/dist/server.js
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal
# Hardening — the dashboard touches gc supervisor over loopback, reads
# .gc/events.jsonl, and runs the `gc` CLI. The `gc bd` write path mutates
# the city's embedded-Dolt bead store directly on disk (under
# %h/gas-city/.beads + %h/gas-city/.gc, taken under a lockfile), so the
# city root must be in ReadWritePaths or those writes fail with EROFS
# despite ProtectHome=read-only. Keep this aligned with GC_CITY_PATH above.
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=%h/.gc %h/gas-city-dashboard %h/gas-city
ProtectHome=read-only
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_UNIX AF_INET
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

[Install]
WantedBy=default.target
